# deploy/supervisor.conf
[program:huawei-cloud-api]
command=/var/www/huawei-cloud-solution-matcher/venv/bin/python -m uvicorn api.main:app --host 127.0.0.1 --port 8000 --workers 4
directory=/var/www/huawei-cloud-solution-matcher
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/huawei-cloud-api.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=PATH="/var/www/huawei-cloud-solution-matcher/venv/bin:%(ENV_PATH)s"
